# docker-compose/php/prod/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777
chown=root:root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

### ini jika pake fpm
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stderr_logfile=/var/log/supervisor/php-fpm.err.log


[program:queue]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
command=php artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker.log
stopwaitsecs=3600


# [program:horizon]
# process_name=%(program_name)s
# directory=/var/www
# command=php artisan horizon
# autostart=true
# autorestart=true
# user=sa
# redirect_stderr=true
# stdout_logfile=/var/log/supervisor/horizon.log
# stopwaitsecs=3600
